# gunicorn.service
[Unit]
Description=Gunicorn application server for Edulytics
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/edulytics
Environment="PATH=/home/ubuntu/edulytics/venv/bin"
ExecStart=/home/ubuntu/edulytics/venv/bin/gunicorn \
    --workers 2 \
    --worker-class sync \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    evaluationWeb.wsgi:application

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
